# memory/openclaw-memory-daily.service
[Unit]
Description=Roland Fleet Memory Daily Maintenance
After=network.target

[Service]
Type=oneshot
User=rolandpg
WorkingDirectory=/home/rolandpg/.openclaw/workspace/memory
ExecStart=/usr/bin/python3 /home/rolandpg/.openclaw/workspace/memory/daily_maintenance.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
