# memory/openclaw-memory-weekly.service
[Unit]
Description=Roland Fleet Memory Weekly Maintenance
After=network.target

[Service]
Type=oneshot
User=rolandpg
WorkingDirectory=/home/rolandpg/.openclaw/workspace/memory
ExecStart=/usr/bin/python3 /home/rolandpg/.openclaw/workspace/memory/weekly_maintenance.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
