# memory/openclaw-memory-weekly.timer
[Unit]
Description=Weekly timer for Roland Fleet Memory maintenance
Requires=openclaw-memory-weekly.service

[Timer]
OnCalendar=Fri *-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
