# memory/openclaw-memory-daily.timer
[Unit]
Description=Daily timer for Roland Fleet Memory maintenance
Requires=openclaw-memory-daily.service

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
