# scripts/zettelforge-rebuild.service
[Unit]
Description=Rebuild ZettelForge Entity Index
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/zettelforge/scripts/rebuild_index.py
WorkingDirectory=/opt/zettelforge
Environment=AMEM_DATA_DIR=/var/lib/zettelforge
User=zettelforge

[Install]
WantedBy=multi-user.target
